# ops/systemd/dr3-trading-agent.service
[Unit]
Description=DR3 IB Trading Agent
After=docker.service network-online.target
Wants=docker.service

[Service]
Type=simple
User=don
Group=don
WorkingDirectory=/home/don/apps/ma-tracker-app/python-service/standalone_agent
ExecStart=/home/don/apps/ma-tracker-app/python-service/.venv/bin/python ib_data_agent.py
Restart=always
RestartSec=2
Nice=-10
# Reserve droplet cores 0-5 for production plane workloads.
CPUAffinity=0 1 2 3 4 5
# Hard cap to avoid unbounded contention; still allows multi-core bursts.
CPUQuota=500%
StandardOutput=append:/home/don/apps/logs/trading-agent.log
StandardError=append:/home/don/apps/logs/trading-agent.log
EnvironmentFile=-/home/don/apps/ma-tracker-app/python-service/.env
EnvironmentFile=/home/don/apps/ma-tracker-app/python-service/standalone_agent/config.env

[Install]
WantedBy=multi-user.target
